# bin/deploy/citytrack-web.conf
# This is an upstart (http://upstart.ubuntu.com/) script
# to run the node.js server on system boot and make it
# manageable with commands such as
# 'start citytrack-web' and 'stop citytrack-web'
#
# This script is to be placed in /etc/init to work with upstart.
#

description "node.js forever server for citytrack-web"
author      "Kostas Stamatoukos <stamatoukosks@gmail.com>"
version "1.0"

expect fork

start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

env HOME=/var/www/citytrack-web

script
# Not sure why $HOME is needed, but it is:
export HOME=$HOME
chdir $HOME
exec /usr/local/bin/node server.js > logs/node.log &
end script

#post-start script
#   # Optionally put a script here that will notifiy you node has (re)started
#end script
